# overlay/usr/lib/systemd/system/hex-wifi.service
[Unit]
Description=Hex Meow Wi-Fi provisioning service
Requires=wpa_supplicant-wlan0.service
Wants=zenohd.service
After=wpa_supplicant-wlan0.service zenohd.service

[Service]
Type=simple
User=root
UMask=0077
RuntimeDirectory=hex-wifi
RuntimeDirectoryMode=0700
ExecStart=/usr/bin/hex-wifi daemon
Restart=on-failure
RestartSec=2
LimitCORE=0
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
