# overlay/usr/lib/systemd/system/zenohd-netjoin.service
# zenohd 的组播加入只在启动时做一次;机器人可能无网线启动,所以 zenohd 本身不等网络
# (本地 client 栈零延迟)。本单元:等 end0 拿到 IPv4(≤180s,无网线则安静退出),若 zenohd
# 的【本次启动实例】曾报 "Cannot join ... multicast"(起得比网络早),try-restart 一次 ——
# 本机 client 会自动重连(实证),LAN 组播自动发现恢复。见 09 §3。
[Unit]
Description=Re-join zenohd multicast scouting once the LAN is up
After=zenohd.service
Wants=zenohd.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'for i in $(seq 1 180); do ip -4 addr show end0 2>/dev/null | grep -q "inet " && break; sleep 1; done; ip -4 addr show end0 2>/dev/null | grep -q "inet " || exit 0; INV=$(systemctl show -p InvocationID --value zenohd); journalctl _SYSTEMD_INVOCATION_ID=$INV --no-pager 2>/dev/null | grep -q "Cannot join" && systemctl try-restart zenohd; exit 0'

[Install]
WantedBy=multi-user.target
