# overlay/usr/lib/systemd/system/serial-getty@ttyGS0.service.d/10-usb-gadget.conf
[Service]
# Avoid systemd's terminal-type probe and terminal reset escape sequences.
# They can be echoed by a Linux cdc_acm host before tio switches the port to
# raw mode, producing an apparent stream of garbage.
Environment=TERM=vt220
ExecStart=
ExecStart=-/sbin/agetty --noclear --wait-cr --issue-file=/etc/issue:/etc/issue.d:/run/issue.d:/usr/lib/issue.d --keep-baud 115200,57600,38400,9600 - ${TERM}
TTYReset=no
